# configs/services/evergreen-device-agent.service
[Unit]
Description=Evergreen Device Enrollment and Policy Agent
Documentation=https://docs.evergreen-os.dev/device-agent
After=network-online.target evergreen-enrollment-greeter.service
Wants=network-online.target
ConditionPathExists=/etc/evergreen/agent/agent.yaml

[Service]
Type=simple
Environment=RUST_LOG=info
ExecStart=/usr/bin/evergreen-device-agent --config /etc/evergreen/agent/agent.yaml
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=6

[Install]
WantedBy=multi-user.target
